# etc/systemd/system/jitsi-sync.service
[Unit]
Description=Synchronize Jitsi certificates
Requisite=docker.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/docker restart jitsi-web