# etc/systemd/system/jitsi-sync.path
Description=Monitor Jitsi certificates

[Path]
PathChanged=/srv/jitsi/certs
Unit=jitsi-sync.service

[Install]
WantedBy=multi-user.target