# etc/systemd/system/nextcloud-sync.service
[Unit]
Description=Synchronize Nextcloud certificates
Requisite=docker.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/docker restart nginx