# etc/systemd/system/ocis-sync.service
[Unit]
Description=Synchronize oCIS certificates
Requisite=docker.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/docker restart ocis