# etc/systemd/system/ocis-sync.path
[Unit]
Description=Monitor oCIS certificates

[Path]
PathChanged=/srv/ocis/certs
Unit=ocis-sync.service

[Install]
WantedBy=multi-user.target